# root/etc/services.d/rclone/run
#!/usr/bin/with-contenv bash

echo "Starting Rclone."

if [ ! -d "${RCLONE_TARGET}" ]; then
  /bin/mkdir -p "${RCLONE_TARGET}"
  chown plex:plex "${RCLONE_TARGET}"
fi

exec s6-setuidgid plex /bin/sh -c "/usr/sbin/rclone --config /config/rclone.conf mount ${RCLONE_MOUNT} ${RCLONE_TARGET} --rc ${RCLONE_ARGS}"