# root/etc/services.d/rclone/run
#!/usr/bin/with-contenv bash

echo "Starting Rclone."

if [ ! -d "/rclone" ]; then
  /bin/mkdir -p "/rclone"
  chown plex:plex "/rclone"
fi

exec s6-setuidgid plex /bin/sh -c "/usr/sbin/rclone --config /config/rclone.conf mount ${RCLONE_MOUNT} ${RCLONE_ARGS}"
